# scripts/com.dotfiles.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Per-user launchd agent: watches the shared dotfiles for uncommitted edits.
     setup.sh templates __HOME__ and installs this to ~/Library/LaunchAgents/. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dotfiles.watch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/Shared/dotfiles/scripts/dotfiles-watch.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/dotfiles-watch.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/dotfiles-watch.err.log</string>
</dict>
</plist>
